# includes.container/usr/lib/systemd/user/vso-tasks-rotation.service
[Unit]
Description=VSO Tasks Rotation

[Service]
ExecStart=/usr/bin/vso tasks rotate

[Install]
WantedBy=multi-user.target
