# includes.container/usr/lib/systemd/system/vanilla-updates.timer
Requires=vanilla-updates.service
After=vanilla-updates.service
[Timer]
OnBootSec=1h
OnUnitActiveSec=1h
Unit=vanilla-updates.service