# includes.container/etc/xdg/autostart/abroot-rollback-notifier.desktop
[Desktop Entry]
Name=ABRoot Rollback Notifier
Exec=/usr/local/bin/abroot-rollback-notifier
Type=Application
