# includes.container/etc/xdg/autostart/vanillaos-vso-tasks-rotator.desktop
[Desktop Entry]
Name=VSO Tasks Rotator
Exec=/usr/bin/vso tasks rotate
Type=Application
